# root/etc/s6-overlay/s6-rc.d/svc-papercut/run
#!/usr/bin/with-contenv bash

exec \
    /etc/init.d/papercut console
